# service/com.ljutaev.invoiceflow.work.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ljutaev.invoiceflow.work</string>

  <!-- Runs the worker continuously: claims jobs and processes them. -->
  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/local-invoice-ocr/.venv/bin/invoiceflow</string>
    <string>work</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>INVOICEFLOW_HOME</key>
    <string>/Users/YOUR_USER/.invoiceflow</string>
    <key>OLLAMA_URL</key>
    <string>http://localhost:11434</string>
    <key>INVOICEFLOW_MODEL</key>
    <string>qwen2.5:14b</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USER/.invoiceflow/worker.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USER/.invoiceflow/worker.err.log</string>
</dict>
</plist>
